# assets/noctalia-polkit.service.in
[Unit]
Description=Noctalia Polkit Authentication Agent
PartOf=graphical-session.target
After=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
ExecStart=@LIBEXECDIR@/noctalia-polkit
Slice=session.slice
TimeoutStopSec=5sec
Restart=on-failure

# Security hardening
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/dbus-1/services
NoNewPrivileges=true
CapabilityBoundingSet=
SystemCallFilter=@system-service
SystemCallArchitectures=native
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=true

[Install]
WantedBy=graphical-session.target
